# debian/atm-tools.atm.service
[Unit]
Documentation=man:atmarpd(8)
Description=Start/stop the atm daemon(s).
DefaultDependencies=no
Before=sysinit.target
Before=networking.service
After=local-fs.target
ConditionCapability=CAP_NET_ADMIN

[Service]
#Environment=OPTIONS=-l syslog
## The EnvironmentFile line only included for backwards/sysvinit compat,
## but using the file to configure is discuraged!
## Instead use 'systemctl edit atm.service' and override the
## Environment or ExecStart lines directly.
#EnvironmentFile=-/etc/default/atm
ExecStart=/usr/sbin/atmarpd $OPTIONS
Restart=on-abort

[Install]
WantedBy=sysinit.target
